# application/third_party/crontab/crontab.txt.all

* * * * * wget -O - http://localhost/index.php/cronjob/minutely_check_for_prayer_time/ >/dev/null 2>&1
